# systemd/ztna-client.service
[Unit]
Description=Zero-Trust Network Access Client
After=network-online.target
Wants=network-online.target
# Allow up to 5 restart attempts within 60 seconds; after that systemd will
# stop retrying and the operator must intervene.  This prevents endless
# restart loops when the config is broken or the controller is unreachable.
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
Type=simple

# The client reads its TOML config from /etc/ztna-client/client.conf.
# No EnvironmentFile needed — the binary handles config loading internally.
ExecStart=/usr/bin/ztna-client serve

WorkingDirectory=/var/lib/ztna-client
# Restart only on non-zero exit (crash/error).  A clean exit via
# `systemctl stop` returns 0 and will NOT trigger a restart.
Restart=on-failure
RestartSec=10
# Give the service 15 seconds to handle SIGTERM and clean up kernel routes
# before systemd sends SIGKILL.
TimeoutStopSec=15

# TUN mode requires CAP_NET_ADMIN for creating/managing the TUN device
# and manipulating ip route/rule entries.  Running as root is the simplest
# valid model for now; a dedicated user + ambient capabilities can be
# added later without changing the unit structure.
#
# NOTE: PrivateDevices is intentionally omitted — TUN mode needs /dev/net/tun.

StateDirectory=ztna-client
RuntimeDirectory=ztna-client
RuntimeDirectoryMode=0755

# --- Hardening (expandable) ---
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
ProtectKernelLogs=yes
ProtectKernelModules=yes
ProtectKernelTunables=yes
ProtectHostname=yes
ProtectClock=yes
ProtectControlGroups=yes
PrivateTmp=yes

RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX AF_NETLINK
RestrictNamespaces=yes
RestrictSUIDSGID=yes
RestrictRealtime=yes
SystemCallArchitectures=native

MemoryDenyWriteExecute=yes
LockPersonality=yes
# DevicePolicy=auto allows access to /dev/net/tun which TUN mode requires.
# PrivateDevices is NOT set for the same reason.
DevicePolicy=auto
DeviceAllow=/dev/net/tun rw
UMask=0077

[Install]
WantedBy=multi-user.target
